# 新的插件商店 Saily/Saily.Daemon/Package/DEBIAN/postinst
echo Installing Saily
echo Version: Major 0.6 Beta b1063

echo Install Daemon Permission
chown root:wheel /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist
chmod 0644 /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist
chmod +x /usr/bin/Saily.Daemon

echo Delete Signal
rm -f /var/root/.SailyBootstrap

rm -rf /var/root/Saily.Daemon # Clean Up

echo Reload Daemon

nohupFile=/usr/bin/nohup
if [ -f "$nohupFile" ]; then
nohup launchctl unload /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist
sleep 1
nohup launchctl load /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist
else
launchctl unload /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist &> /dev/null
launchctl load /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist &> /dev/null
fi

echo Install Application Permission
chmod 0755 /Applications/Saily.app/Saily
chown -R root:wheel /Applications/Saily.app
chown -R root:wheel /Applications/Saily.app/*

# Uicache
# uicache /Applications/Saily.app
# Broken uicache on Chimera, sorry to say.
# echo 'Sorry, we need to do uicache -a cause it is broken by: ________'
# uicache /Applications/Saily.app
# su mobile -c uicache
# Damemon will run uicache -a

# Force dpkg to continue
# echo "com.Lakr233.Saily install" | dpkg --set-selections  
exit 0
